# docker/supervisor/conf.d/watcher.conf
[program:watcher]
command=php watcher.php
autostart=false
autorestart=false
startretries=0
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/etc/supervisor/log/watcher.log
